# s6/etc/s6-overlay/s6-rc.d/litestream-init/up
#!/command/with-contenv
/usr/bin/bash /etc/s6-overlay/s6-rc.d/litestream-init/init