# s6/etc/s6-overlay/s6-rc.d/wordpress/run
#!/command/with-contenv sh
cd /app
export Wordpress__Constants__FQDBDIR=/tmp/
/usr/bin/dotnet /app/app.dll