# s6/etc/s6-overlay/s6-rc.d/litestream/run
#!/command/with-contenv bash
echo "Syncing SQLite WAL"

if [[ $BUCKET_LOCATION ]]; then
  /opt/bin/litestream replicate -o /app/wordpress/wp-content/database/.ht.sqlite $BUCKET_LOCATION
else
  echo "BUCKET_LOCATION env var not set, skipping litestream sync"
  /usr/bin/sleep 99d
fi

