# s6/etc/s6-overlay/s6-rc.d/litestream/run
#!/command/with-contenv bash
echo "Syncing SQLite WAL"
/usr/bin/sleep 99d